# lumiera-whatsapp.service
[Unit]
Description=Lumiera WhatsApp Copilot API
After=network.target

[Service]
Type=simple
User=ceeai
Group=ceeai
WorkingDirectory=/home/ceeai/whatsapp_api
Environment="PATH=/home/ceeai/whatsapp_api/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ceeai/whatsapp_api"
ExecStart=/home/ceeai/whatsapp_api/venv/bin/python -m uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

# Logging: Use journald (view with: journalctl -u lumiera-whatsapp -f)
# Application logs to logs/app.log with automatic rotation
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
